# packaging/privacy-manager-server.service

[Unit]
Description=privacy manager server

[Service]
ExecStart=/usr/bin/privacy-manager-server

[Install]
WantedBy=multi-user.target
